# consul/hello-consul/hello-consul_1.0-1/etc/init/hello-consul.conf
description "Hello Consul"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

respawn

exec hello-consul
