# google/stackdriver_monitoring/config/prometheus/node_exporter.conf
start on filesystem or runlevel [2345]
exec /usr/bin/node_exporter
