# consul/upstart/consul.server.conf
description "Consul server"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

respawn

setuid consul
setgid consul

exec consul agent -server -config-dir /etc/consul.d/server -client 0.0.0.0
