# fiat/fiat-web/etc/init/fiat.conf
description "fiat"

start on filesystem or runlevel [2345]

expect fork

stop on stopping spinnaker

exec sudo -u spinnaker -g spinnaker /opt/fiat/bin/fiat 2>&1 >> /var/log/spinnaker/fiat/fiat.log &
